# jobs/macos/Users/mr-fancy/Library/LaunchAgents/com.chrisrjones.exp-launchd-notif-center.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple/DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.chrisrjones.exp-launchd-notif-center</string>
    <key>ProgramArguments</key>
    <array>
      <string>sh</string>
      <string>-c</string>
      <string>$HOME/.local/bin/mac-notif-center-test.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/launchd/com.chrisrjones.test.macos-notif-center-test.err</string>

    <key>StandardOutPath</key>
    <string>/usr/local/var/log/launchd/com.chrisrjones.test.macos-notif-center-test.out</string>
  </dict>
</plist>
