# jobs/macOS/Users/mr-fancy/Library/LaunchAgents/com.chrisrjones.ln_bus_pirate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple/DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.chrisrjones.ln_bus_pirate.plist</string>

    <key>ProgramArguments</key>
    <array>
      <string>sh</string>
      <string>-c</string>
      <string>/Users/capin/.local/bin/macos-bp-ln.sh</string>
    </array>

    <key>WatchPaths</key>
    <array>
      <string>/dev/tty.usbserial-A700dzDq</string>
    </array>

    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/ln_bus_pirate.err</string>

    <key>StandardOutPath</key>
    <string>/usr/local/var/log/ln_bus_pirate.out</string>
  </dict>
</plist>



