# jobs/macOS/Users/mr-fancy/Library/LaunchAgents/com.chrisrjones.brew-update-install-list.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <!-- # NOTE: the below key has not been implemented --> 
    <!-- <key>RootDirectory</key> -->
    <!-- <string>~</string> -->
    <!-- <key>GroupName</key> -->
    <!-- <string>staff</string> -->
    <!-- <key>UserName</key> -->
    <!-- <string>capin</string> -->
    <!-- <key>WorkingDirectory</key> -->
    <!-- <string>~</string> -->
    <!-- <key>EnableGlobbing</key> -->
    <!-- <true/> -->

    <key>Label</key>
    <string>com.chrisrjones.brew-update-install-list</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/usr/local/bin:/usr/bin:/bin</string>
      <!-- <key>HOME</key> -->
      <!-- <string>/Users/capin</string> -->
    </dict>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/sh</string>
      <string>-c</string>
      <string>launchctl setenv USER "${USER}"</string>
      <string>launchctl setenv HOME "${HOME}"</string>
    </array>

    <key>ProgramArguments</key>
    <!-- <key>Program</key> -->
    <array>
      <string>/bin/sh</string>
      <string>-c</string>
      <string>~/.local/bin/brew-update-install-list.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
      <key>SuccessfulExit</key>
      <false/>
    </dict>

    <key>StandardErrorPath</key>
    <string>/Users/capin/logs/brew-update-install-list.err</string>
    <!-- <string>$HOME/logs/brew-update-install-list.err</string> -->
    <!-- <string>/Users/$USER$/logs/brew-update-install-list.err</string> -->
    <!-- <string>~/logs/brew-update-install-list.err</string> -->

    <key>StandardOutPath</key>
    <string>/Users/capin/logs/brew-update-install-list.out</string>
    <!-- <string>$HOME/logs/brew-update-install-list.out</string> -->
    <!-- <string>/Users/$USER$/logs/brew-update-install-list.out</string> -->
    <!-- <string>~/logs/brew-update-install-list.out</string> -->

    <!-- <key>WorkingDirectory</key> -->
    <!-- <string>$HOME</string> -->
    <!-- <string>~</string> -->

    <!-- <key>StartCalendarInterval</key> -->
    <!-- <array> -->
      <!-- 	<dict> -->
        <!-- 		<key>Hour</key> -->
        <!-- 		<integer>23</integer> -->
        <!-- 		<key>Minute</key> -->
        <!-- 		<integer>0</integer> -->
        <!-- 		<key>Month</key> -->
        <!-- 		<integer>0</integer> -->
        <!-- 		<key>Weekday</key> -->
        <!-- 		<integer>4</integer> -->
        <!-- 	</dict> -->
      <!-- 	<dict> -->
        <!-- 		<key>Hour</key> -->
        <!-- 		<integer>23</integer> -->
        <!-- 		<key>Minute</key> -->
        <!-- 		<integer>0</integer> -->
        <!-- 		<key>Month</key> -->
        <!-- 		<integer>1</integer> -->
        <!-- 		<key>Weekday</key> -->
        <!-- 		<integer>4</integer> -->
        <!-- 	</dict> -->
      <!-- 	<dict> -->
        <!-- 		<key>Hour</key> -->
        <!-- 		<integer>23</integer> -->
        <!-- 		<key>Minute</key> -->
        <!-- 		<integer>0</integer> -->
        <!-- 		<key>Month</key> -->
        <!-- 		<integer>2</integer> -->
        <!-- 		<key>Weekday</key> -->
        <!-- 		<integer>4</integer> -->
        <!-- 	</dict> -->
      <!-- 	<dict> -->
        <!-- 		<key>Hour</key> -->
        <!-- 		<integer>23</integer> -->
        <!-- 		<key>Minute</key> -->
        <!-- 		<integer>0</integer> -->
        <!-- 		<key>Month</key> -->
        <!-- 		<integer>3</integer> -->
        <!-- 		<key>Weekday</key> -->
        <!-- 		<integer>4</integer> -->
        <!-- 	</dict> -->
      <!-- 	<dict> -->
        <!-- 		<key>Hour</key> -->
        <!-- 		<integer>23</integer> -->
        <!-- 		<key>Minute</key> -->
        <!-- 		<integer>0</integer> -->
        <!-- 		<key>Month</key> -->
        <!-- 		<integer>4</integer> -->
        <!-- 		<key>Weekday</key> -->
        <!-- 		<integer>4</integer> -->
        <!-- 	</dict> -->
      <!-- 	<dict> -->
        <!-- 		<key>Hour</key> -->
        <!-- 		<integer>23</integer> -->
        <!-- 		<key>Minute</key> -->
        <!-- 		<integer>0</integer> -->
        <!-- 		<key>Month</key> -->
        <!-- 		<integer>5</integer> -->
        <!-- 		<key>Weekday</key> -->
        <!-- 		<integer>4</integer> -->
        <!-- 	</dict> -->
      <!-- 	<dict> -->
        <!-- 		<key>Hour</key> -->
        <!-- 		<integer>23</integer> -->
        <!-- 		<key>Minute</key> -->
        <!-- 		<integer>0</integer> -->
        <!-- 		<key>Month</key> -->
        <!-- 		<integer>6</integer> -->
        <!-- 		<key>Weekday</key> -->
        <!-- 		<integer>4</integer> -->
        <!-- 	</dict> -->
      <!-- 	<dict> -->
        <!-- 		<key>Hour</key> -->
        <!-- 		<integer>23</integer> -->
        <!-- 		<key>Minute</key> -->
        <!-- 		<integer>0</integer> -->
        <!-- 		<key>Month</key> -->
        <!-- 		<integer>7</integer> -->
        <!-- 		<key>Weekday</key> -->
        <!-- 		<integer>4</integer> -->
        <!-- 	</dict> -->
      <!-- 	<dict> -->
        <!-- 		<key>Hour</key> -->
        <!-- 		<integer>23</integer> -->
        <!-- 		<key>Minute</key> -->
        <!-- 		<integer>0</integer> -->
        <!-- 		<key>Month</key> -->
        <!-- 		<integer>8</integer> -->
        <!-- 		<key>Weekday</key> -->
        <!-- 		<integer>4</integer> -->
        <!-- 	</dict> -->
      <!-- 	<dict> -->
        <!-- 		<key>Hour</key> -->
        <!-- 		<integer>23</integer> -->
        <!-- 		<key>Minute</key> -->
        <!-- 		<integer>0</integer> -->
        <!-- 		<key>Month</key> -->
        <!-- 		<integer>9</integer> -->
        <!-- 		<key>Weekday</key> -->
        <!-- 		<integer>4</integer> -->
        <!-- 	</dict> -->
      <!-- 	<dict> -->
        <!-- 		<key>Hour</key> -->
        <!-- 		<integer>23</integer> -->
        <!-- 		<key>Minute</key> -->
        <!-- 		<integer>0</integer> -->
        <!-- 		<key>Month</key> -->
        <!-- 		<integer>10</integer> -->
        <!-- 		<key>Weekday</key> -->
        <!-- 		<integer>4</integer> -->
        <!-- 	</dict> -->
      <!-- 	<dict> -->
        <!-- 		<key>Hour</key> -->
        <!-- 		<integer>23</integer> -->
        <!-- 		<key>Minute</key> -->
        <!-- 		<integer>0</integer> -->
        <!-- 		<key>Month</key> -->
        <!-- 		<integer>11</integer> -->
        <!-- 		<key>Weekday</key> -->
        <!-- 		<integer>4</integer> -->
        <!-- 	</dict> -->
      <!-- </array> -->
  </dict>
</plist>
