# jobs/macOS/Users/mr-fancy/Library/LaunchAgents/com.chrisrjones.brew-update-Brewfile.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.chrisrjones.brew-update-Brewfile</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/sh</string>
      <string>-c</string>
      <string>$HOME/.local/bin/brew-update-Brewfile.sh</string>
    </array>

    <!-- <key>RunAtLoad</key> -->
    <!-- <true/> -->

    <!-- <key>KeepAlive</key> -->
    <!-- <dict> -->
    <!--   <key>SuccessfulExit</key> -->
    <!--   <false/> -->
    <!-- </dict> -->

    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/launchd/com.chrisrjones.brew-update-Brewfile.err</string>

    <key>StandardOutPath</key>
    <string>/usr/local/var/log/launchd/com.chrisrjones.brew-update-Brewfile.out</string>

    <key>StartCalendarInterval</key>
    <array>
      	<dict>
        		<key>Hour</key>
        		<integer>23</integer>
        		<key>Minute</key>
        		<integer>0</integer>
        		<key>Month</key>
        		<integer>0</integer>
        		<key>Weekday</key>
        		<integer>4</integer>
        	</dict>
      	<dict>
        		<key>Hour</key>
        		<integer>23</integer>
        		<key>Minute</key>
        		<integer>0</integer>
        		<key>Month</key>
        		<integer>1</integer>
        		<key>Weekday</key>
        		<integer>4</integer>
        	</dict>
      	<dict>
        		<key>Hour</key>
        		<integer>23</integer>
        		<key>Minute</key>
        		<integer>0</integer>
        		<key>Month</key>
        		<integer>2</integer>
        		<key>Weekday</key>
        		<integer>4</integer>
        	</dict>
      	<dict>
        		<key>Hour</key>
        		<integer>23</integer>
        		<key>Minute</key>
        		<integer>0</integer>
        		<key>Month</key>
        		<integer>3</integer>
        		<key>Weekday</key>
        		<integer>4</integer>
        	</dict>
      	<dict>
        		<key>Hour</key>
        		<integer>23</integer>
        		<key>Minute</key>
        		<integer>0</integer>
        		<key>Month</key>
        		<integer>4</integer>
        		<key>Weekday</key>
        		<integer>4</integer>
        	</dict>
      	<dict>
        		<key>Hour</key>
        		<integer>23</integer>
        		<key>Minute</key>
        		<integer>0</integer>
        		<key>Month</key>
        		<integer>5</integer>
        		<key>Weekday</key>
        		<integer>4</integer>
        	</dict>
      	<dict>
        		<key>Hour</key>
        		<integer>23</integer>
        		<key>Minute</key>
        		<integer>0</integer>
        		<key>Month</key>
        		<integer>6</integer>
        		<key>Weekday</key>
        		<integer>4</integer>
        	</dict>
      	<dict>
        		<key>Hour</key>
        		<integer>23</integer>
        		<key>Minute</key>
        		<integer>0</integer>
        		<key>Month</key>
        		<integer>7</integer>
        		<key>Weekday</key>
        		<integer>4</integer>
        	</dict>
      	<dict>
        		<key>Hour</key>
        		<integer>23</integer>
        		<key>Minute</key>
        		<integer>0</integer>
        		<key>Month</key>
        		<integer>8</integer>
        		<key>Weekday</key>
        		<integer>4</integer>
        	</dict>
      	<dict>
        		<key>Hour</key>
        		<integer>23</integer>
        		<key>Minute</key>
        		<integer>0</integer>
        		<key>Month</key>
        		<integer>9</integer>
        		<key>Weekday</key>
        		<integer>4</integer>
        	</dict>
      	<dict>
        		<key>Hour</key>
        		<integer>23</integer>
        		<key>Minute</key>
        		<integer>0</integer>
        		<key>Month</key>
        		<integer>10</integer>
        		<key>Weekday</key>
        		<integer>4</integer>
        	</dict>
      	<dict>
        		<key>Hour</key>
        		<integer>23</integer>
        		<key>Minute</key>
        		<integer>0</integer>
        		<key>Month</key>
        		<integer>11</integer>
        		<key>Weekday</key>
        		<integer>4</integer>
        	</dict>
      </array>
  </dict>
</plist>
