# jobs/macos/Users/mr-fancy/Library/LaunchAgents/com.chrisrjones.brew-updatedb.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple/DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.chrisrjones.brew-updatedb</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/sh</string>
      <string>-c</string>
      <string>$HOME/.local/bin/brew-updatedb.sh</string>
    </array>

    <!-- <key>RunAtLoad</key> -->
    <!-- <true/> -->

    <!-- <key>KeepAlive</key> -->
    <!-- <true/> -->
    <dict>
      <key>SuccessfulExit</key>
      <false/>
    </dict>

    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/launchd/com.chrisrjones-brew-updatedb.err</string>

    <key>StandardOutPath</key>
    <string>/usr/local/var/log/launchd/com.chrisrjones-brew-updatedb.out</string>

    <key>StartCalendarInterval</key>
    <array>
      <dict>
        <key>Day</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>11</integer>
      </dict>
    </array>
  </dict>
</plist>
