# jobs/macos/Users/mr-fancy/Library/LaunchAgents/com.chrisrjones.rm-ds-store-files.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple/DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.chrisrjones.rm-ds-store-files</string>

    <key>ProgramArguments</key>
    <array>
      <string>sh</string>
      <string>-c</string>
      <string>$HOME/.local/bin/rm-ds-store-files.sh</string>
    </array>

    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/launchd/com.chrisrjones.rm_ds_store_files.err</string>

    <key>StandardOutPath</key>
    <string>/usr/local/var/log/launchd/com.chrisrjones.rm_ds_store_files.out</string>

    <!-- <key>RunAtLoad</key> -->
    <!-- <true/> -->

    <key>StartCalendarInterval</key>
    <array>
      <dict>
        <key>Day</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>2</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>3</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>4</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>5</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>6</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>7</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>8</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>9</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>10</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>11</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>12</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>13</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>14</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>15</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>16</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>17</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>18</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>19</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>20</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>21</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>22</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>23</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>24</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>25</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>26</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>27</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>28</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>29</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>30</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>31</integer>
        <key>Hour</key>
        <integer>5</integer>
      </dict>
    </array>
  </dict>
</plist>

