# jobs/macOS/Users/mr-fancy/Library/LaunchAgents/com.chrisrjones.brew-update-install-list.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <!-- # NOTE: the below key has not been implemented --> 
    <!-- <string>~</string> -->
    <!-- <key>RootDirectory</key> -->
    <!-- <key>GroupName</key> -->
    <!-- <string>staff</string> -->
    <!-- <key>UserName</key> -->
    <!-- <string>capin</string> -->
    <!-- <key>WorkingDirectory</key> -->
    <!-- <string>~</string> -->
    <key>Label</key>
    <string>com.chrisrjones.brew-update-install-list</string>
    <key>Program</key>
    <!-- <string>/bin/sh</string> -->
    <!-- <string>-c</string> -->
    <string>/Users/capin/.local/bin/brew-update-install-list.sh</string>
    
    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/bin:/usr/bin:/usr/local/bin</string>
    </dict>

    <key>StandardErrorPath</key>
    <string>/Users/capin/logs/brew-update-install-list.err</string>

    <key>StandardOutPath</key>
    <string>/Users/capin/logs/brew-update-install-list.out</string>

    <key>WorkingDirectory</key>
    <string>/Users/capin</string>

	<!-- <key>StartCalendarInterval</key> -->
	<!-- <array> -->
	<!-- 	<dict> -->
	<!-- 		<key>Hour</key> -->
	<!-- 		<integer>23</integer> -->
	<!-- 		<key>Minute</key> -->
	<!-- 		<integer>0</integer> -->
	<!-- 		<key>Month</key> -->
	<!-- 		<integer>0</integer> -->
	<!-- 		<key>Weekday</key> -->
	<!-- 		<integer>4</integer> -->
	<!-- 	</dict> -->
	<!-- 	<dict> -->
	<!-- 		<key>Hour</key> -->
	<!-- 		<integer>23</integer> -->
	<!-- 		<key>Minute</key> -->
	<!-- 		<integer>0</integer> -->
	<!-- 		<key>Month</key> -->
	<!-- 		<integer>1</integer> -->
	<!-- 		<key>Weekday</key> -->
	<!-- 		<integer>4</integer> -->
	<!-- 	</dict> -->
	<!-- 	<dict> -->
	<!-- 		<key>Hour</key> -->
	<!-- 		<integer>23</integer> -->
	<!-- 		<key>Minute</key> -->
	<!-- 		<integer>0</integer> -->
	<!-- 		<key>Month</key> -->
	<!-- 		<integer>2</integer> -->
	<!-- 		<key>Weekday</key> -->
	<!-- 		<integer>4</integer> -->
	<!-- 	</dict> -->
	<!-- 	<dict> -->
	<!-- 		<key>Hour</key> -->
	<!-- 		<integer>23</integer> -->
	<!-- 		<key>Minute</key> -->
	<!-- 		<integer>0</integer> -->
	<!-- 		<key>Month</key> -->
	<!-- 		<integer>3</integer> -->
	<!-- 		<key>Weekday</key> -->
	<!-- 		<integer>4</integer> -->
	<!-- 	</dict> -->
	<!-- 	<dict> -->
	<!-- 		<key>Hour</key> -->
	<!-- 		<integer>23</integer> -->
	<!-- 		<key>Minute</key> -->
	<!-- 		<integer>0</integer> -->
	<!-- 		<key>Month</key> -->
	<!-- 		<integer>4</integer> -->
	<!-- 		<key>Weekday</key> -->
	<!-- 		<integer>4</integer> -->
	<!-- 	</dict> -->
	<!-- 	<dict> -->
	<!-- 		<key>Hour</key> -->
	<!-- 		<integer>23</integer> -->
	<!-- 		<key>Minute</key> -->
	<!-- 		<integer>0</integer> -->
	<!-- 		<key>Month</key> -->
	<!-- 		<integer>5</integer> -->
	<!-- 		<key>Weekday</key> -->
	<!-- 		<integer>4</integer> -->
	<!-- 	</dict> -->
	<!-- 	<dict> -->
	<!-- 		<key>Hour</key> -->
	<!-- 		<integer>23</integer> -->
	<!-- 		<key>Minute</key> -->
	<!-- 		<integer>0</integer> -->
	<!-- 		<key>Month</key> -->
	<!-- 		<integer>6</integer> -->
	<!-- 		<key>Weekday</key> -->
	<!-- 		<integer>4</integer> -->
	<!-- 	</dict> -->
	<!-- 	<dict> -->
	<!-- 		<key>Hour</key> -->
	<!-- 		<integer>23</integer> -->
	<!-- 		<key>Minute</key> -->
	<!-- 		<integer>0</integer> -->
	<!-- 		<key>Month</key> -->
	<!-- 		<integer>7</integer> -->
	<!-- 		<key>Weekday</key> -->
	<!-- 		<integer>4</integer> -->
	<!-- 	</dict> -->
	<!-- 	<dict> -->
	<!-- 		<key>Hour</key> -->
	<!-- 		<integer>23</integer> -->
	<!-- 		<key>Minute</key> -->
	<!-- 		<integer>0</integer> -->
	<!-- 		<key>Month</key> -->
	<!-- 		<integer>8</integer> -->
	<!-- 		<key>Weekday</key> -->
	<!-- 		<integer>4</integer> -->
	<!-- 	</dict> -->
	<!-- 	<dict> -->
	<!-- 		<key>Hour</key> -->
	<!-- 		<integer>23</integer> -->
	<!-- 		<key>Minute</key> -->
	<!-- 		<integer>0</integer> -->
	<!-- 		<key>Month</key> -->
	<!-- 		<integer>9</integer> -->
	<!-- 		<key>Weekday</key> -->
	<!-- 		<integer>4</integer> -->
	<!-- 	</dict> -->
	<!-- 	<dict> -->
	<!-- 		<key>Hour</key> -->
	<!-- 		<integer>23</integer> -->
	<!-- 		<key>Minute</key> -->
	<!-- 		<integer>0</integer> -->
	<!-- 		<key>Month</key> -->
	<!-- 		<integer>10</integer> -->
	<!-- 		<key>Weekday</key> -->
	<!-- 		<integer>4</integer> -->
	<!-- 	</dict> -->
	<!-- 	<dict> -->
	<!-- 		<key>Hour</key> -->
	<!-- 		<integer>23</integer> -->
	<!-- 		<key>Minute</key> -->
	<!-- 		<integer>0</integer> -->
	<!-- 		<key>Month</key> -->
	<!-- 		<integer>11</integer> -->
	<!-- 		<key>Weekday</key> -->
	<!-- 		<integer>4</integer> -->
	<!-- 	</dict> -->
	<!-- </array> -->
</dict>
</plist>
