# jobs/Linux/arch/etc/profile.d/capin_hb_paths.sh
#!/bin/bash

if [ "$USER" = "capin" ]; then
  export PATH=$HOME/hombrew/sbin:$HOME/homebrew/bin:$PATH
fi
