# config/brew/macOS/Sierra/Library/LaunchDaemons/com.chrisrjones.sshd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<!-- <key>Disabled</key> -->
	<!-- <true/> -->
	<key>Label</key>
	<string>com.chrisrjones.sshd</string>
  <!-- OpenSSH must be run as super user -->
	<!-- <key>InitGroups</key>
	<true/> -->
	<key>Program</key>
	<!-- <string>/usr/libexec/sshd-keygen-wrapper</string> -->
	<string>/usr/local/sbin/sshd</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/sbin/sshd</string>
		<!--	<string>-i</string> -->
    <!-- <string>-d</string> --> <!-- useful for debugging --> 
		<string>-f</string>
		<string>/usr/local/etc/ssh/sshd_config</string>
		<string>-E</string>
		<string>/usr/local/var/log/sshd.log</string>
		<string>-o</string>
		<string>PidFile=/usr/local/var/run/ssh/sshd.pid</string>
		<string>-D</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>Sockets</key>
	<dict>
		<key>Listeners</key>
		<dict>
			<key>SockServiceName</key>
			<string>ssh</string>
			<key>Bonjour</key>
			<array>
				<string>ssh</string>
				<string>sftp-ssh</string>
			</array>
		</dict>
	</dict>
	<!-- <key>inetdCompatibility</key> 
	<dict>
		<key>Wait</key>
		<false/>
		<key>Instances</key>
		<integer>42</integer>
	</dict> -->
	<key>StandardErrorPath</key>
	<string>/usr/local/var/log/sshd-error.log</string>
	<key>StandardOutPath</key>
	<string>/usr/local/var/log/sshd-out.log</string>
	<key>SHAuthorizationRight</key>
	<string>system.preferences</string>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
</dict>
</plist>
